# systemd/grid-bot-sol.service
[Unit]
Description=Grid Trading Bot - SOL Monitor
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/binance-trading-bot-v3
ExecStart=/usr/bin/node /root/binance-trading-bot-v3/grid-bot-cli-v5.mjs monitor --name live-sol-bot
Restart=always
RestartSec=10
StandardOutput=append:/root/binance-trading-bot-v3/logs/live-sol-bot.log
StandardError=append:/root/binance-trading-bot-v3/logs/live-sol-bot.log
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
