# systemd/grid-bot-btc.service
[Unit]
Description=Grid Trading Bot - BTC Monitor
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/binance-trading-bot-v3
ExecStart=/usr/bin/node /root/binance-trading-bot-v3/grid-bot-cli-v5.mjs monitor --name live-btc-bot
Restart=always
RestartSec=10
StandardOutput=append:/root/binance-trading-bot-v3/logs/live-btc-bot.log
StandardError=append:/root/binance-trading-bot-v3/logs/live-btc-bot.log
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
